# meta-gateway-tora/recipes-gateway/gateway-watch-config/gateway-watch-config/confwatch.path
[Path]
Unit=confwatch.service
PathModified=/etc/gateway/config

[Install]
WantedBy=multi-user.target
