# meta-gateway-tora/recipes-gateway/gateway-watch-wifi/gateway-watch-wifi/netwatch.path
[Path]
Unit=netwatch.service
PathModified=/etc/gateway/network

[Install]
WantedBy=multi-user.target
